# crontab
* * * * * web python /code/etl/main.py >> /var/log/cron.log 2>&1
